# docker/s6-rc.d/bridge/run
#!/command/execlineb -P
fdmove -c 2 1
cd /app
/usr/local/bin/python main.py
